# init/ubuntu/openflow_rcs.service
[Unit]
Description=OpenFlow related RCs
After=syslog.target

[Service]
Type=simple

[Install]
WantedBy=multi-user.target
